# pymailadmin.service
[Unit]
Description=PyMailAdmin - Web administration interface for mail server
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/gunicorn \
    --worker-class sync \
    --workers 2 \
    --bind 127.0.0.1:8686 \
    --worker-tmp-dir /tmp \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --error-logfile /var/log/pymailadmin/error.log \
    --timeout 120 \
    --preload app:app
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
TimeoutStopSec=10

# Logs
StandardOutput=null
StandardError=null

# Unprivileged user
User=pymailadmin
Group=pymailadmin

# Namespace
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true

# Capacities
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Files
ReadWritePaths=/var/log/pymailadmin
ReadOnlyPaths=/var/www/pymailadmin
InaccessiblePaths=/etc/passwd

# Environment
WorkingDirectory=/var/www/pymailadmin
Environment=PYTHONUNBUFFERED=1

# Restart policy
Restart=always
RestartSec=10

# Watchdog
WatchdogSec=30

[Install]
WantedBy=multi-user.target
